# home/Library/LaunchAgents/net.serubin.sleepwatcher.plist.tmpl
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>net.serubin.sleepwatcher</string>
	<key>ProgramArguments</key>
	<array>
		<string>{{ output "brew" "--prefix" | trim }}/opt/sleepwatcher/sbin/sleepwatcher</string>
		<string>-V</string>
		<string>-s</string>
		<string>{{ .chezmoi.homeDir }}/.config/sleepwatcher/sleep</string>
		<string>-w</string>
		<string>{{ .chezmoi.homeDir }}/.config/sleepwatcher/wakeup</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>KeepAlive</key>
	<true/>
	<key>ProcessType</key>
	<string>Background</string>
</dict>
</plist>
